# squid_proxy_manager/rootfs/etc/services.d/squid-proxy-manager/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: Squid Proxy Manager
# Runs the main Python application
# ==============================================================================

# Write to log file for debugging (in case stdout isn't captured)
LOGFILE="/data/squid_proxy_manager/startup.log"
mkdir -p /data/squid_proxy_manager
chown -R 1000:1000 /data/squid_proxy_manager || true
chmod 750 /data/squid_proxy_manager || true
touch "$LOGFILE"
chown 1000:1000 "$LOGFILE" || true

log() {
    local timestamp=$(date '+%Y-%m-%d %H:%M:%S')
    echo "[$timestamp] $1" >> "$LOGFILE"
    echo "[$timestamp] $1"
}

log "=========================================="
log "Squid Proxy Manager - Service Starting"
log "=========================================="
log "PWD: $(pwd)"
log "USER: $(whoami)"
log "Listing /app:"
ls -la /app/ >> "$LOGFILE" 2>&1 || log "Failed to list /app"

# Check Python
log "Checking Python..."
if command -v python3 &> /dev/null; then
    PYTHON_VERSION=$(python3 --version 2>&1)
    log "Python found: ${PYTHON_VERSION}"
else
    log "ERROR: Python3 not found!"
    exit 1
fi

# Check main.py
if [ ! -f "/app/main.py" ]; then
    log "ERROR: /app/main.py not found!"
    exit 1
fi
log "main.py found"

# Set environment
export PYTHONUNBUFFERED=1
export PYTHONPATH=/app

# Get log level from config
if command -v bashio &> /dev/null 2>&1; then
    export LOG_LEVEL=$(bashio::config 'log_level' 'info' 2>/dev/null || echo 'info')
else
    export LOG_LEVEL='info'
fi
log "LOG_LEVEL: ${LOG_LEVEL}"

# Build Squid Docker image if it doesn't exist
# (No longer needed in process-based architecture, keeping for reference)
# log "Checking for Squid proxy Docker image..."
# if [ -f "/app/build_squid_image.sh" ]; then
#     if /app/build_squid_image.sh >> "$LOGFILE" 2>&1; then
#         log "✓ Squid proxy image ready"
#     else
#         log "WARNING: Failed to build Squid proxy image. Check $LOGFILE for details."
#     fi
# else
#     log "WARNING: build_squid_image.sh not found, skipping image build"
# fi

log "Starting Python application..."
log "Command: exec python3 /app/main.py"

# Execute Python
exec s6-setuidgid 1000:1000 python3 /app/main.py
